# systemd/figurehue-healthcheck.service
[Unit]
Description=FigureHue health check and recovery

[Service]
Type=oneshot
ExecStart=INSTALL_DIR/figurehue_healthcheck.sh
